# scripts/com.jollyfarmer.tradeshow-email-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jollyfarmer.tradeshow-email-ingest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ernest/Documents/GitHub/JF/Trade Shows/scripts/email_ingest.py</string>
        <string>--watch</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ernest/Documents/GitHub/JF/Trade Shows</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ernest/Documents/GitHub/JF/Trade Shows/inbox/email-ingest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ernest/Documents/GitHub/JF/Trade Shows/inbox/email-ingest.log</string>
</dict>
</plist>
